# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D -e

[program:fail2ban]
command=/usr/bin/fail2ban-server -xf start